# data/lolbins-cti.json
[{"source_project":"LOLBins-CTI","tool_name":"BITSAdmin.exe","platform":"windows","description":"BITSAdmin is a command line tool used to create and manage BITS Jobs.","native_purpose":"BITSAdmin is a command line tool used to create and manage BITS Jobs.; Paths: c:\\windows\\system32\\bitsadmin.exe, c:\\windows\\syswow64\\bitsadmin.exe","known_behaviors":[],"command_examples":[],"indicators":["c:\\windows\\system32\\bitsadmin.exe","c:\\windows\\syswow64\\bitsadmin.exe"],"similar_tools":[],"categories":["TA0009:Persistence:T1197","TA0005:Defense Evasion:T1197","TA0008:Lateral Movement:T1570","TA0010:Exfiltration:T1048.003","TA0011:Command and Control:T1105","privilege:User"],"references":["https://github.com/CTI-Driven/LOLBins"]},{"source_project":"LOLBins-CTI","tool_name":"Certutil.exe","platform":"windows","description":"Command-line utility that can be used to obtain certificate authority information and configure Certificate Services.","native_purpose":"Command-line utility that can be used to obtain certificate authority information and configure Certificate Services.; Paths: c:\\windows\\system32\\certutil.exe, c:\\windows\\syswow64\\certutil.exe","known_behaviors":[],"command_examples":[],"indicators":["c:\\windows\\system32\\certutil.exe","c:\\windows\\syswow64\\certutil.exe"],"similar_tools":[],"categories":["TA0009:Collection:T1560.001","TA0005:Defense Evasion:T1140","TA0011:Command and Control:T1105","TA0005:Defense Evasion:T1553.004 ","privilege:User"],"references":["https://github.com/CTI-Driven/LOLBins"]},{"source_project":"LOLBins-CTI","tool_name":"PsExec[.]exe","platform":"windows","description":"PsExec is a free Microsoft tool that can be used to execute a program on another computer. It is used by IT administrators and attackers.","native_purpose":"PsExec is a free Microsoft tool that can be used to execute a program on another computer. It is used by IT administrators and attackers.; Paths: C:\\SysinternalsSuite\\PsExec[.]exe, C:\\.\\PsExec[.]exe","known_behaviors":[],"command_examples":[],"indicators":["C:\\SysinternalsSuite\\PsExec[.]exe","C:\\.\\PsExec[.]exe"],"similar_tools":[],"categories":["TA0002:Execution:T1569.002","TA0003:Persistence:T1136.002","TA0004:Privilege Escalation:T1543.003","TA0008:Lateral Movement:T1570","TA0008:Lateral Movement:T1021.002","privilege:User"],"references":["https://github.com/CTI-Driven/LOLBins"]},{"source_project":"LOLBins-CTI","tool_name":"Reg.exe","platform":"windows","description":"Reg is a Windows utility used to interact with the Windows Registry. It can be used at the command-line interface to query, add, modify, and remove information. [1]\\n\\nUtilities such as Reg are known to be used by persistent threats.","native_purpose":"Reg is a Windows utility used to interact with the Windows Registry. It can be used at the command-line interface to query, add, modify, and remove information. [1]\\n\\nUtilities such as Reg are known to be used by persistent threats.; Paths: c:\\windows\\system32\\reg.exe, c:\\windows\\syswow64\\reg.exe","known_behaviors":[],"command_examples":[],"indicators":["c:\\windows\\system32\\reg.exe","c:\\windows\\syswow64\\reg.exe"],"similar_tools":[],"categories":["TA0005:Defense Evasion:T1112","TA0007:Discovery:T1012","TA0006:Credential Access:T1552.002","privilege:User","privilege:Administrator"],"references":["https://github.com/CTI-Driven/LOLBins"]},{"source_project":"LOLBins-CTI","tool_name":"Regsvr32.exe","platform":"windows","description":"Adversaries may abuse Regsvr32.exe to proxy execution of malicious code. Regsvr32.exe is a command-line program used to register and unregister object linking and embedding controls, including dynamic link libraries (DLLs), on Windows systems. The Regsvr32.exe binary may also be signed by Microsoft. (Citation: Microsoft Regsvr32)\n\nMalicious usage of Regsvr32.exe may avoid triggering security tools that may not monitor execution of, and modules loaded by, the regsvr32.exe process because of allowlists or false positives from Windows using regsvr32.exe for normal operations. Regsvr32.exe can also be used to specifically bypass application control using functionality to load COM scriptlets to execute DLLs under user permissions. Since Regsvr32.exe is network and proxy aware, the scripts can be loaded by passing a uniform resource locator (URL) to file on an external Web server as an argument during invocation. This method makes no changes to the Registry as the COM object is not actually registered, only executed. (Citation: LOLBAS Regsvr32) This variation of the technique is often referred to as a \"Squiblydoo\" and has been used in campaigns targeting governments. (Citation: Carbon Black Squiblydoo Apr 2016) (Citation: FireEye Regsvr32 Targeting Mongolian Gov)\n\nRegsvr32.exe can also be leveraged to register a COM Object used to establish persistence via [Component Object Model Hijacking](https://attack.mitre.org/techniques/T1546/015). (Citation: Carbon Black Squiblydoo Apr 2016)","native_purpose":"Adversaries may abuse Regsvr32.exe to proxy execution of malicious code. Regsvr32.exe is a command-line program used to register and unregister object linking and embedding controls, including dynamic link libraries (DLLs), on Windows systems. The Regsvr32.exe binary may also be signed by Microsoft. (Citation: Microsoft Regsvr32)\n\nMalicious usage of Regsvr32.exe may avoid triggering security tools that may not monitor execution of, and modules loaded by, the regsvr32.exe process because of allowlists or false positives from Windows using regsvr32.exe for normal operations. Regsvr32.exe can also be used to specifically bypass application control using functionality to load COM scriptlets to execute DLLs under user permissions. Since Regsvr32.exe is network and proxy aware, the scripts can be loaded by passing a uniform resource locator (URL) to file on an external Web server as an argument during invocation. This method makes no changes to the Registry as the COM object is not actually registered, only executed. (Citation: LOLBAS Regsvr32) This variation of the technique is often referred to as a \"Squiblydoo\" and has been used in campaigns targeting governments. (Citation: Carbon Black Squiblydoo Apr 2016) (Citation: FireEye Regsvr32 Targeting Mongolian Gov)\n\nRegsvr32.exe can also be leveraged to register a COM Object used to establish persistence via [Component Object Model Hijacking](https://attack.mitre.org/techniques/T1546/015). (Citation: Carbon Black Squiblydoo Apr 2016); Paths: c:\\windows\\system32\\regsvr32.exe, c:\\windows\\syswow64\\regsvr32.exe","known_behaviors":[],"command_examples":[],"indicators":["c:\\windows\\system32\\regsvr32.exe","c:\\windows\\syswow64\\regsvr32.exe"],"similar_tools":[],"categories":["TA0005:Defense Evasion:T1218.010","privilege:User"],"references":["https://github.com/CTI-Driven/LOLBins"]},{"source_project":"LOLBins-CTI","tool_name":"rundll32.exe","platform":"windows","description":"The rundll32.exe program can be called to execute an arbitrary binary. Adversaries may take advantage of this functionality to proxy execution of code to avoid triggering security tools that may not monitor execution of the rundll32.exe process because of whitelists or false positives from Windows using rundll32.exe for normal operations.","native_purpose":"The rundll32.exe program can be called to execute an arbitrary binary. Adversaries may take advantage of this functionality to proxy execution of code to avoid triggering security tools that may not monitor execution of the rundll32.exe process because of whitelists or false positives from Windows using rundll32.exe for normal operations.; Paths: c:\\windows\\system32\\rundll32.exe, c:\\windows\\syswow64\\rundll32.exe","known_behaviors":[],"command_examples":[],"indicators":["c:\\windows\\system32\\rundll32.exe","c:\\windows\\syswow64\\rundll32.exe"],"similar_tools":[],"categories":["TA0004:Privilege Escalation:T1548.002","TA0006:Persis Credential Accesstence:T1555.004,T1003.001","TA0005:Defense Evasion:T1218.011","privilege:User"],"references":["https://github.com/CTI-Driven/LOLBins"]},{"source_project":"LOLBins-CTI","tool_name":"schtasks[.]exe","platform":"windows","description":"Adversaries may abuse the Windows Task Scheduler to perform task scheduling for initial or recurring execution of malicious code.","native_purpose":"Adversaries may abuse the Windows Task Scheduler to perform task scheduling for initial or recurring execution of malicious code.; Paths: c:\\windows\\system32\\schtasks.exe, c:\\windows\\syswow64\\schtasks.exe","known_behaviors":[],"command_examples":[],"indicators":["c:\\windows\\system32\\schtasks.exe","c:\\windows\\syswow64\\schtasks.exe"],"similar_tools":[],"categories":["TA0002:Execution:T1053.005","TA0003:Persistence:T1053.005","TA0004:Privilege Escalation:T1053.005","privilege:User"],"references":["https://github.com/CTI-Driven/LOLBins"]}]